# dev/vagrant/dolibarrdev/puphpet/puppet/modules/mailcatcher/templates/etc/init/mailcatcher.conf.erb
# mailcatcher - mock smtp server
#
# mailCatcher runs a super simple SMTP server which catches any
# message sent to it to display in a web interface.

description     "mock smtp server"

start on startup
stop on shutdown

setuid mailcatcher
setgid mailcatcher

script
 exec nohup <%= @mailcatcher_path %>/mailcatcher<%= @options.join(' ') %> -f  >> /var/log/mailcatcher/mailcatcher.log 2>&1
end script
